# ansible/roles/nuxt-supervisord/templates/supervisord.conf.j2
[program:{{ server_name }}]
directory={{ nuxt_dir }}
command=node ./node_modules/.bin/nuxt start
user=www-data
autostart=true
autorestart=true
