# ansible/roles/nuxt-supervisord/templates/supervisord.conf.j2
[program:{{ server_name }}]
directory={{ root_dir }}
command=node ./node_modules/.bin/nuxt start
user=ubuntu
autostart=true
autorestart=true
